# templates/glances.service.j2
# {{ ansible_managed }}
#
[Unit]
Description=Glances
After=syslog.target network.target auditd.service

[Service]
ExecStart=/usr/bin/glances -w
Restart=on-abort
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target